# release/src/router/dnsmasq/debian/postrm
#!/bin/sh
set -e

if [ purge = "$1" ]; then
   update-rc.d dnsmasq remove >/dev/null
fi
