# release/src/router/avahi-0.6.31/initscript/gentoo/avahi-daemon.in
#!/sbin/runscript
# Copyright 1999-2007 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

opts="reload"

depend() {
	before netmount nfsmount
	use net
	need dbus
}

start() {
	ebegin "Starting avahi-daemon"
	@sbindir@/avahi-daemon -D
	eend $?
}

stop() {
	ebegin "Stopping avahi-daemon"
	@sbindir@/avahi-daemon -k
	eend $?
}

reload() {
	ebegin "Reloading avahi-daemon"
	@sbindir@/avahi-daemon -r
	eend $?
}
