# release/src/router/avahi-0.6.31/initscript/gentoo/avahi-dnsconfd.in
#!/sbin/runscript
# Copyright 1999-2007 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

opts="reload"

depend() {
	need avahi-daemon
	use net
}

start() {
	ebegin "Starting avahi-dnsconfd"
	@sbindir@/avahi-dnsconfd -D
	eend $?
}

stop() {
	ebegin "Stopping avahi-dnsconfd"
	@sbindir@/avahi-dnsconfd -k
	eend $?
}

reload() {
	ebegin "Reloading avahi-dnsconfd"
	@sbindir@/avahi-dnsconfd -r
	eend $?
}
