# release/src/router/usb-modeswitch/usb-modeswitch-upstart.conf
start on usb-modeswitch-upstart
task
script
	exec /usr/sbin/usb_modeswitch_dispatcher --switch-upstart $UMS_PARAM
end script
