# release/src/router/rp-pppoe/scripts/pppoe-init.in
#!/bin/sh
#
# pppoe                     This script starts or stops a PPPoE connection
#
# chkconfig: 2345 99 01
# description: Connects to PPPoE provider
#
# LIC: GPL
#
# Copyright (C) 2000 Roaring Penguin Software Inc.  This software may
# be distributed under the terms of the GNU General Public License, version
# 2 or any later version.

# Source function library if it exists
test -r /etc/rc.d/init.d/functions && . /etc/rc.d/init.d/functions

# From AUTOCONF
prefix=@prefix@
exec_prefix=@exec_prefix@

# Paths to programs
START=@sbindir@/pppoe-start
STOP=@sbindir@/pppoe-stop
STATUS=@sbindir@/pppoe-status
case "$1" in
    start)
        echo -n "Bringing up PPPoE link"

	$START
	if [ $? = 0 ] ; then
		touch /var/lock/subsys/pppoe
	        echo_success
	else
		echo_failure
	fi
        echo ""
        ;;

    stop)
        echo -n "Shutting down PPPoE link"

	$STOP > /dev/null 2>&1
	if [ $? = 0 ] ; then
		rm -f /var/lock/subsys/pppoe
	        echo_success
	else
		echo_failure
	fi
        echo ""
        ;;

    restart)
	$0 stop
	$0 start
	;;

    status)
	$STATUS
	;;

    *)
        echo "Usage: pppoe {start|stop|restart|status}"
        exit 1
esac

exit 0
