# release/src/router/mt-daapd-svn-1696/contrib/init.d/mtdaapd-bsd.templ
#!/bin/sh
#
# PROVIDE: mtdaapd mt-daapd
# REQUIRE: DAEMON
# KEYWORD: shutdown
#
# mt-daapd daemon.
#
# Created by Patrick Lemmens for NetBSD and FreeBSD.

. /etc/rc.subr

name="mtdaapd"
realname="mt-daapd"
rcvar=$name
command="@prefix@/bin/${realname}"
pidfile="/var/run/${realname}.pid"
sig_stop="SIGINT"
 
load_rc_config $name
run_rc_command "$1"