# release/src/router/usb-modeswitch-2.1.1/usb_modeswitch@.service
[Unit]
Description=USB_ModeSwitch

[Service]
Type=oneshot
ExecStart=/usr/sbin/usb_modeswitch_dispatcher --switch-systemd %I
